# bin/puma.service
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
Type=forking

User=project
Environment="BEGOOD_PATH=/home/project/Backend"
Environment="BEGOOD_VERSION=v2"

WorkingDirectory=/home/project/Backend/lib/sinatra
ExecStart=/usr/local/bin/puma -C /home/project/Backend/lib/sinatra/puma.rb
PIDFile=/home/project/Backend/log/puma.pid
Restart=no

[Install]
WantedBy=multi-user.target
